# units/lazy_mount.path
[Unit]
Description=Path activator of lazy mount service

[Path]
PathExists=/tmp/.lazy_mount
