# units/tizen-fstrim-user.timer
[Unit]
Description=Timer for tizen-fstrim

[Timer]
OnCalendar=*-*-* 06:00:00
